# renewd/renewd.service
[Unit]
Description=Kerberos Credentials Automatic Rewewal

[Service]
Type=forking
ExecStart=/usr/sbin/renewd

[Install]
WantedBy=multi-user.target

